# root/etc/uci-defaults/99_homeproxy_ui
#!/bin/sh

# Create symbolic link for HomeProxy dashboard
[ -d /etc/homeproxy/ui ] && [ ! -L /www/hp ] && ln -s /etc/homeproxy/ui /www/hp

exit 0
